# deploy/tdl-stream-bot.service
[Unit]
Description=TDL DingTalk Stream Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/bots/tdl/backend
EnvironmentFile=/opt/bots/tdl/backend/.env
Environment=PYTHONPATH=/opt/bots/tdl/backend
ExecStart=/opt/bots/tdl/backend/.venv/bin/python -m app.integrations.dingtalk_stream_bot
Restart=always
RestartSec=3
User=tdl

[Install]
WantedBy=multi-user.target